# openrc/nodexmastree
#!/sbin/openrc-run

name="nodexmastree"
command="/usr/bin/node"
command_args="/home/igor/nodexmastree/main.js"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
